# supervisord.conf
[supervisord]
nodaemon=true

[program:gogs]
command=/var/gogs/gogs/gogs web
environment=HOME="/var/gogs",USER="gogs"
