# files/app/gogs/docker/s6/gogs/run
#!/bin/sh

export USER=git
exec gosu $USER /opt/gogs/gogs web
