# s6/gogs/run
#!/bin/sh

export USER=gogs
exec gosu $USER /app/gogs/gogs web
