# s6/gogs/run
#!/bin/sh

export USER=git
exec gosu $USER /app/gogs/gogs web
